# frontend/scripts/kaspi-dump.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd-агент Каспи-демпинг-воркера Croon. Держит воркер живым: запускает при логине,
  перезапускает если упал (KeepAlive), переживает сон/перезагрузку Mac.

  Плейсхолдеры __WORKDIR__ / __NODE__ / __SECRET__ / __LOGINUSER__ / __LOGINPASS__ /
  __MERCHANT__ / __SITE__ подставляет `kaspi-dump.sh install` при установке (значения
  берутся из frontend/.env). Готовый plist кладётся в ~/Library/LaunchAgents/ —
  В РЕПОЗИТОРИЙ НЕ КОММИТИТСЯ (содержит секреты).

  Воркер сам .env НЕ читает — поэтому переменные передаём здесь, в EnvironmentVariables.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>kz.croon.kaspi-dump</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>scripts/kaspi-cabinet-worker.mjs</string>
  </array>

  <!-- Рабочий каталог = frontend/ (воркер ищет scripts/ и .kaspi-session.json относительно него) -->
  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SITE</key>
    <string>__SITE__</string>
    <key>DUMPING_SECRET</key>
    <string>__SECRET__</string>
    <key>KASPI_LOGIN_USER</key>
    <string>__LOGINUSER__</string>
    <key>KASPI_LOGIN_PASS</key>
    <string>__LOGINPASS__</string>
    <key>MERCHANT_UID</key>
    <string>__MERCHANT__</string>
    <key>LOOP_MIN</key>
    <string>2</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Стартовать при загрузке агента (вход в систему) -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Перезапускать при падении (в т.ч. когда воркер сам выходит из-за протухшей
       сессии — тогда в логе будет видно «нужен повторный --login»). Throttle 60с,
       чтобы при постоянной ошибке не крутить рестарты в цикле. -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/tmp/croon-kaspi-dump.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/croon-kaspi-dump.log</string>

  <!-- Процессорный приоритет фоновой задачи -->
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
